# CONFIG_FILES/original-files/4.5.0.19/%usr%lib%systemd%user%mms-engine.service
[Unit]
Description=MMS engine service

[Service]
ExecStart=/usr/bin/invoker --type=generic /usr/sbin/mms-engine -o syslog
